# traefik-addon/rootfs/etc/services.d/traefik/run
#!/command/with-contenv bashio
echo "=== DEBUG === Listing /data:"
ls -l /data

echo "=== DEBUG === Starting Traefik..."
exec traefik --configFile=/data/traefik/traefik.yml

