# traefik-addon/rootfs/etc/services.d/traefik/run
#!/command/with-contenv bashio
echo "Starting Traefik..."
exec traefik --configFile=/data/traefik/traefik.yml

