# traefik-addon/rootfs/etc/services.d/traefik/run
#!/command/with-contenv bashio
LEGO_CA_CERTIFICATES=/data/traefik/certs/root_ca.crt exec traefik --configFile=/data/traefik/traefik.yml

