# traefik-addon/rootfs/etc/services.d/traefik/run
#!/usr/bin/execlineb -p
exec traefik --configFile=/data/traefik/traefik.yml

